# install/restart_daemon.crontab
@reboot /var/www/mycodo/mycodo/scripts/restart_daemon.py --continuouscheck &
