# install/restart_daemon.crontab
@reboot /var/www/mycodo/env_py3/bin/python3 /var/www/mycodo/mycodo/scripts/restart_daemon.py --continuouscheck &
